# systemd/dosnap-autoclean.timer
[Unit]
Description=Daily Cleanup of dosnap Snapshots

[Timer]
OnBootSec=10m
OnUnitActiveSec=1d

[Install]
WantedBy=timers.target
