# systemd/dosnap-autoclean.timer
[Unit]
Description=Daily Cleanup of dosnap Snapshots

[Timer]
OnCalendar=daily
RandomizedDelaySec=1min

[Install]
WantedBy=timers.target
