# systemd/dosnap-create.timer
[Unit]
Description=Auto Snapshot for dosnap

[Timer]
OnCalendar=01/2:00
OnStartupSec=2m
RandomizedDelaySec=30

[Install]
WantedBy=timers.target
