# systemd/dosnap-create.service
[Unit]
Description=Auto Snapshot for dosnap

[Service]
Type=simple
ExecStart=/usr/bin/dosnap create --auto --all
